# autostart/albert.desktop
[Desktop Entry]
Type=Application
Name=Albert
Exec=albert
Icon=albert
Comment=Alfred clone for Linux
