# autostart/copyq.desktop
[Desktop Entry]
Type=Application
Name=CopyQ
Exec=copyq
Icon=copyq
Comment=Copy/Paste manager for Linux
